# smartyfi-remote/rootfs/etc/services.d/zerotier/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: ZeroTier One
# Runs ZeroTier One
# ==============================================================================
declare -a options

bashio::log.info "Starting SmartyFi Remote ..."

# Note sure what this does.
options+=(-U)
options+=("-p$(bashio::addon.port 9993)")

#Run Bridge Script
exec /etc/services.d/zerotier/bridge.sh &


# Run the ZeroTier One
exec zerotier-one "${options[@]}"


